# .crontab
25 2 * * * tarsnap-backup >> /home/noel/crontab.log
