# .crontab
25 2 * * * clamavscanner.sh >> ~/clamav.log
